# Prometheus/Prometheus/Servicios/node_export.service
[Unit]
Description=Node Export service

[Service]
User=root
ExecStart=/opt/node_exporter/node_exporter

[Install]
WantedBy=multi-user.target
